# packaging/linux/npu-webhook.service
[Unit]
Description=npu-webhook 个人知识库服务
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/npu-webhook
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
